# deploy/sub2api.bt-supervisor.example.ini
[program:sub2api]
directory=/www/wwwroot/s2.ai80.vip
command=/www/wwwroot/s2.ai80.vip/run-sub2api.sh
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=15
stopsignal=TERM
stdout_logfile=/www/wwwlogs/sub2api.out.log
stderr_logfile=/www/wwwlogs/sub2api.err.log
redirect_stderr=false

# 由 deploy/deploy-zero-downtime.sh 每次部署时自动生成 run-sub2api.sh，
# 里面会显式 export GIN_MODE / SERVER_HOST / SERVER_PORT / DATA_DIR。
# 这样即使系统里存在 SERVER_PORT=8787 之类的旧环境变量，也不会误继承。
